# opsworks/clickhouse/templates/clickhouse.service.erb
[Unit]
Description=ClickHouse Server (analytic DBMS for big data)
Requires=network-online.target
After=network-online.target

[Service]
Type=simple
User=<%= node['clickhouse']['user'] %>
Group=<%= node['clickhouse']['group'] %>
Restart=always
RestartSec=30
RuntimeDirectory=clickhouse-server
ExecStart=/usr/bin/clickhouse-server --config=<%= node['clickhouse']['conf_path'] %>/config.xml --pid-file=<%= node['clickhouse']['run_path'] %>/clickhouse-server.pid
LimitCORE=infinity
LimitNOFILE=500000
LimitNPROC=40960
CapabilityBoundingSet=CAP_NET_ADMIN CAP_IPC_LOCK CAP_SYS_NICE

[Install]
WantedBy=multi-user.target